# crontab
33 9 * * 1-5 /Users/charliegetzen/.local/bin/brew_manager &> /tmp/brew.log && /usr/local/bin/terminal-notifier -title "Updated brew" -message "/tmp/brew.log"
