# crontab
33 9 * * 1-5 /Users/charliegetzen/github.com/cgetzen/dotfiles/brew_manager.sh &> /tmp/brew.log && /usr/local/bin/terminal-notifier -title "Updated brew" -message "/tmp/brew.log"
